# dotfilesd/src/daemon_files/dotfiles-Daemon.service
[Service]
ExecStart=/usr/local/bin/dotfiles daemon
Restart=always
RestartSec=5
