# dotfilesd/src/daemon_files/dotfiles-daemon.service
[Service]
ExecStart=/usr/local/bin/dotfilesd daemon
Restart=always
RestartSec=5
