# bundle/linux/usr/lib/systemd/user/fig.service
[Unit]
Description=Start Fig
After=multi-user.target
Wants=graphical.target

[Service]
ExecStart=/usr/bin/fig_desktop

